# scripts/launchd/com.kidan.local-staging.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kidan.local-staging</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/opt/php@8.2/bin/php</string>
    <string>artisan</string>
    <string>serve</string>
    <string>--host=127.0.0.1</string>
    <string>--port=8000</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/media/Desktop/Kidan App</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>5</integer>

  <key>StandardOutPath</key>
  <string>/tmp/kidan-php.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/kidan-php.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>LARAVEL_SAIL</key>
    <string>0</string>
  </dict>
</dict>
</plist>
